# device/supervisord.conf
[supervisord]
nodaemon=true

[program:monitoring]
command=/app/device/start_monitoring.sh
autostart=true
autorestart=true
